# systemd/container-watcher.service
[Unit]
.
.

[Service]
Type=notify
ExecStart=
NotifyAccess=all
Restart=always
.
.

[Install]
WantedBy= 
