# systemd/go-doc-6534.socket
[Unit]
Description=On Demand offline godocs

[Socket]
ListenStream=6534
#Accept=true
[Install]
WantedBy=sockets.target
