# systemd/ssh-socks-5564.socket
[Unit]
Description=On Demand Socks proxy into Work

[Socket]
ListenStream=9119
#Accept=true
[Install]
WantedBy=sockets.target
